# system_files/etc/systemd/system/incus.service.d/10-env.conf
[Service]
Environment=INCUS_UI=/usr/lib/opt/incus/ui/
